# display_monitor.service
[Unit]
Description=Display Monitor Service
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Face-rec-main
Environment="PATH=/home/pi/face-rec/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/pi/face-rec/bin/python -u /home/pi/Face-rec-main/display_monitor.py
Restart=always
RestartSec=10
StandardOutput=append:/home/pi/Face-rec-main/display_monitor.log
StandardError=append:/home/pi/Face-rec-main/display_monitor.log

[Install]
WantedBy=multi-user.target 